# supervisor_conf/entity.conf
[program:Entity_Extraction]
command=/usr/bin/python2.7 /home/minu/TextAnalytics/entity/entity.py
autostart=true
autorestart=true